# 01.相关资料/shijiu/CoreTest.ini
[program:CoreTest]
command=dotnet alcWeb.dll
directory=/home/wwwroot
environment=ASPNETCORE__ENVIRONMENT=Production
user=root
stopsignal=INT
autostart=true
autorestart=true
startsecs=3
stderr_logfile=/var/log/ossoffical.err.log 
stdout_logfile=/var/log/ossoffical.out.log